# artifacts/speaker-agent.service
[Unit]
Description=Bluetooth speaker agent service
Documentation=https://github.com/japoch/sonos-streaming/bluetooth-receiver
After=multi-user.target

[Service]
Type=simple
Environment=PYTHONUNBUFFERED=1
ExecStart=/usr/bin/python /home/sonos/ws/gh-sonos-streaming/bluetooth-receiver/speaker-agent.py

[Install]
WantedBy=multi-user.target
